# cron/crontab
* * * * * /script.sh >> /var/log/cron.log