# cron/crontab
0 3 * * * /script.sh >> /var/log/cron.log