# cron/crontab
0 0 * * * /script.sh >> /var/log/cron.log